# rtorrent/include/supervisord.ini
[supervisord]
nodaemon=true

[program:rtorrent]
priority=2
redirect_stderr=true
startsecs=5
command=rtorrent